# backend/unifiedviews-backend-shared/src/deb/control/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


set_java() {
 JAVA_VERSION=`java -version 2>&1 | grep "java version" | awk '{print $3}' | tr -d \" | awk '{split($0, array, ".")} END{print array[2]}'`
 echo -n "check java version"
 if [ ! $JAVA_VERSION -eq 7 ]; then
     JAVA_NAME=`update-java-alternatives  -l | grep 7 | cut -d' ' -f1`
     if [ -n $JAVA_NAME ] ; then
         /usr/sbin/update-java-alternatives --jre-headless --set $JAVA_NAME
     echo " - Java version is set to 1.7"
     else
        echo " ERROR: java 7 is not installed"
     fi
 else
    echo  " - ok"
 fi
}


prepare_dir() { 
    mkdir -p /var/log/unifiedviews/backend
    chown tomcat7:tomcat7 /var/log/unifiedviews/backend
    
    mkdir -p /usr/share/unifiedviews
    chown tomcat7:tomcat7 /usr/share/unifiedviews

    mkdir -p /var/log/unifiedviews/backend
    chown tomcat7:tomcat7 /var/log/unifiedviews /var/log/unifiedviews/backend

    mkdir -p /var/lib/unifiedviews/backend/working
    chown tomcat7:tomcat7  /var/lib/unifiedviews/backend  /var/lib/unifiedviews/backend/working

    mkdir -p /var/lib/unifiedviews/target
    chown tomcat7:tomcat7 /var/lib/unifiedviews/target /var/lib/unifiedviews/target

    mkdir -p /var/lib/unifiedviews/target/
    chown tomcat7:tomcat7 /var/lib/unifiedviews/backend/working
    
    mkdir -p ~tomcat7/.eclipse
    chown -R tomcat7:tomcat7 ~tomcat7/.eclipse

    chown tomcat7:tomcat7 /usr/share/unifiedviews
    chmod g+w /usr/share/unifiedviews
    chmod +x  /usr/sbin/run_unifiedviews_backend
    chown tomcat7:tomcat7 /usr/sbin/run_unifiedviews_backend
  
    chmod +x  /etc/init.d/unifiedviews-backend

    chown tomcat7:tomcat7 /usr/share/tomcat7
    
    mkdir -p /var/lib/unifiedviews/target/lib
    chown tomcat7:tomcat7 /var/lib/unifiedviews/target/lib
    
    mkdir -p /var/lib/unifiedviews/target/dpu
    chown tomcat7:tomcat7 /var/lib/unifiedviews/target/dpu
    
    if [ -e "/var/www/dump/" ]; then
        chown -R tomcat7:tomcat7 /var/www/dump/ 
    fi
    
}

service_backend_start() {
    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d unifiedviews-backend start || true
    else
        service unifiedviews-backend start || true
    fi
}

case "$1" in
  configure)
    set_java
    prepare_dir
    service_backend_start
  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
  ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    
  ;;
esac


# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
